# etc/systemd/system/scbackup.timer
[Unit]
Description=Daily scbackup timer

[Timer]
OnCalendar=*-*-* 1:00:00
Persistent=true

[Install]
WantedBy=basic.target
